# scripts/deploy/hartonomous-model-ingestion.service
[Unit]
Description=Hartonomous Model Ingestion
After=network.target

[Service]
Type=notify
User=hartonomous
WorkingDirectory=/srv/www/hartonomous/model-ingestion
ExecStart=/usr/bin/dotnet /srv/www/hartonomous/model-ingestion/ModelIngestion.dll
Restart=always
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=hartonomous-model-ingestion
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false

[Install]
WantedBy=multi-user.target
