# deploy/hartonomous-neo4j-sync.service
[Unit]
Description=Hartonomous Neo4j Sync
After=network.target neo4j.service

[Service]
Type=notify
User=hartonomous
WorkingDirectory=/srv/www/hartonomous/neo4j-sync
ExecStart=/usr/bin/dotnet /srv/www/hartonomous/neo4j-sync/Neo4jSync.dll
Restart=always
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=hartonomous-neo4j-sync
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false

[Install]
WantedBy=multi-user.target
