# deploy/hartonomous-api.service
[Unit]
Description=Hartonomous API
After=network.target

[Service]
Type=notify
User=hartonomous
WorkingDirectory=/srv/www/hartonomous/api
ExecStart=/usr/bin/dotnet /srv/www/hartonomous/api/Hartonomous.Api.dll
Restart=always
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=hartonomous-api
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false

[Install]
WantedBy=multi-user.target
