# deploy/hartonomous-ces-consumer.service
[Unit]
Description=Hartonomous CES Consumer
After=network.target

[Service]
Type=notify
User=hartonomous
WorkingDirectory=/srv/www/hartonomous/ces-consumer
ExecStart=/usr/bin/dotnet /srv/www/hartonomous/ces-consumer/CesConsumer.dll
Restart=always
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=hartonomous-ces-consumer
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false

[Install]
WantedBy=multi-user.target
